# vm/chef/cookbooks/meanstack/files/conf/supervisor-mean.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=redirect_stdout
loglevel=info
user=root

[program:api]
directory=/sites/api
command=/bin/bash -c "source $NVM_DIR/nvm.sh && node index.js"
autorestart=true
autostart=true
priority=1

[program:web]
directory=/sites/web
command=/bin/bash -c "source $NVM_DIR/nvm.sh && source /sites/api/env_vars.sh && ng serve --host 0.0.0.0 --port %(ENV_WEB_PORT)s"
autostart=true
autorestart=true
priority=2

[program:webserver]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=4
