# vm/chef/cookbooks/resourcespace/files/resourcespace.cron
#!/bin/sh
wget -q -r http://localhost/pages/tools/cron_copy_hitcount.php
