# vm/chef/cookbooks/orocrm/files/supervisor.conf

[program:oro_web_socket]
command=php /opt/orocrm/${ORO_PROJECT}/bin/console gos:websocket:server --env=prod
numprocs=1
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stderr_logfile=/var/log/oroscript.err.log
stdout_logfile=/var/log/oroscript.out.log

[program:oro_message_consumer]
command=php /opt/orocrm/${ORO_PROJECT}/bin/console oro:message-queue:consume --env=prod
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stderr_logfile=/var/log/oroscript.err.log
stdout_logfile=/var/log/oroscript.out.log

