# vm/chef/cookbooks/sonarqube/templates/sonar.service.erb
[Unit]
Description=SonarQube service
After=syslog.target network.target

[Service]
Type=forking

LimitNOFILE=65536
LimitNPROC=4096

ExecStart=/opt/sonarqube/bin/linux-x86-64/sonar.sh start
ExecStop=/opt/sonarqube/bin/linux-x86-64/sonar.sh stop

User=sonar
Group=sonar
Restart=always

[Install]
WantedBy=multi-user.target
