# vm/chef/cookbooks/elasticsearch/files/certshare.service
[Unit]
Description=ELK Certshare Service

[Service]
WorkingDirectory=/tmp/certs/files
ExecStart=/usr/bin/python3 -m http.server
