# docker/spark/3/debian11/3.4/supervisor.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=redirect_stdout
loglevel=info
user=root

[program:spark_master]
command=/opt/spark/bin/spark-class org.apache.spark.deploy.master.Master --host "%(ENV_SPARK_MASTER_HOST)s" --port "%(ENV_SPARK_MASTER_PORT)s" --webui-port "%(ENV_SPARK_MASTER_WEBUI_PORT)s"
autorestart=true
autostart=%(ENV_START_MASTER)s
priority=1
# Redirect log to the main process
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:spark_worker]
command=/opt/spark/bin/spark-class org.apache.spark.deploy.worker.Worker --webui-port "%(ENV_SPARK_WORKER_WEBUI_PORT)s" "%(ENV_SPARK_MASTER)s"
autorestart=true
autostart=%(ENV_START_WORKER)s
priority=1
# Redirect log to the main process
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:spark_history]
command=/opt/spark/bin/spark-class org.apache.spark.deploy.history.HistoryServer
autorestart=true
autostart=%(ENV_START_HISTORY)s
priority=2
# Redirect log to the main process
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
